# backend/systemd/realestate-scraper-metro.service
[Unit]
Description=Real Estate Scraper - Metro Scope
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/aksai/projects/realestate/backend
ExecStart=/usr/bin/python3 etl_extract_v3.py --scope metro --limit 5000
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
